# opt/backup-to-remote/backup-to-usb-@.service
[Unit]
Description=Backup to USB Flash Disk
BindsTo=dev-%i.device

[Service]
Type=simple
ExecStart=/opt/backup-to-remote/backup-to-remote.sh -l usb
