# opt/backup-to-remote/backup-to-remote.cron
# Copie vers Rclone tous les dimanches à 4h du matin
0 3 * * sun root /opt/backup-to-remote/backup-to-remote.sh rclone -q