# opt/mysqlbackup/mysqlbackup.cron
# Fichier de configuration du cron pour MysqlBackup /etc/cron.d/mysqlbackup
# Exécution toutes les lundi soir à 3h00
0 3 * * 1 root /opt/mysqlbackup/mysqlbackup.sh >> /var/log/mysqlbackup.log
