# opt/mysqlbackup/mysqlbackup.cron
#!/bin/sh
# /etc/cron.weekly/mysqlbackup

echo "`date +'%Y-%m-%d %H:%M:%S'` Démarrage script de sauvegarde"
/opt/mysqlbackup/mysqlbackup.sh >> /var/log/mysqlbackup.log

