# opt/mysqlbackup/mysqlbackup.cron
#!/bin/bash
# /etc/cron.weekly/mysqlbackup

SCRIPT_DIR=$( dirname "$(readlink -f "$0")" )
SCRIPT_FILE="$SCRIPT_DIR/mysqlbackup.sh"

$SCRIPT_FILE >> /var/log/mysqlbackup.log
