# opt/rclonebackup/rclonebackup.cron
#!/bin/sh
# /etc/cron.daily/rclonebackup

echo "`date +'%Y-%m-%d %H:%M:%S'` Démarrage script de sauvegarde"
/opt/rclonebackup/rclonebackup.sh >> /var/log/rclonebackup.log

