# opt/rclonebackup/rclonebackup.cron
#!/bin/bash
# /etc/cron.weekly/rclonebackup

SCRIPT_DIR=$( dirname "$(readlink -f "$0")" )
SCRIPT_FILE="$SCRIPT_DIR/rclonebackup.sh"

$SCRIPT_FILE >> /var/log/rclonebackup.log
