# opt/rclonebackup/rclonebackup.cron
# Fichier de configuration du cron pour RcloneBackup /etc/cron.d/rclonebackup
# Exécution toutes les mardi soir à 3h00
0 3 * * 2 root /opt/rclonebackup/rclonebackup.sh -vl >> /var/log/rclonebackup.log 2&>1
