# etc/profile.d/dircolors.sh
#!/bin/sh
# Active la coloration correcte des liens symboliques
eval `dircolors -b`

