# etc/update-motd.d/colors.sh
#Colours
red="\e[00;31m"
RED="\e[01;31m"

green="\e[00;32m"
GREEN="\e[01;32m"

brown="\e[00;33m"
YELLOW="\e[01;33m"

blue="\e[00;34m"
BLUE="\e[01;34m"

purple="\e[00;35m"
PURPLE="\e[01;35m"

cyan="\e[00;36m"
CYAN="\e[01;36m"

white="\e[00;37m"
WHITE="\e[01;37m"

NC="\e[00m"

